# deploy/roles/osp/templates/supervisor/osp-workers.conf.j2

[program:osp-workers]
command={{ osp_src }}/env/bin/rqworker
process_name=%(program_name)s_%(process_num)02d
numprocs=2
