# deploy/roles/osp/templates/supervisor/rq-dashboard.conf.j2

[program:rq-dashboard]
command={{ osp_src }}/env/bin/uwsgi --ini {{ osp_uwsgi_config }}/rq-dashboard.ini
