# deploy/roles/osp/templates/supervisor/tika.conf.j2

[program:tika]
command=jara -jar {{ tika_src }}
