# deploy/roles/osp/templates/supervisor/tika-server.conf.j2

[program:tika-server]
command=java -jar {{ tika_path }}/{{ tika_jar }}
autostart=false
