# deploy/roles/workers/templates/supervisor/es-test.conf.j2

[program:es-test]
command=/usr/share/elasticsearch/bin/elasticsearch --config=/etc/elasticsearch/elasticsearch.test.yml
autostart=false
