# deploy/roles/osp/templates/osp.conf.j2

[program:osp-workers]
command={{ osp_dest }}/env/bin/rqworker
process_name=%(program_name)s_%(process_num)02d
numprocs={{ osp_num_workers }}
