# deploy/roles/osp/templates/rq-dashboard.conf.j2

[program:rq-dashboard]
command={{ osp_src }}/env/bin/python3 {{ osp_src }}/osp/web/rq_dash.py
