# deploy/roles/osp/templates/rq.conf.j2

[program:rq-dashboard]
command={{ osp_dest }}/env/bin/python3 {{ osp_dest }}/bin/rq_dash
